# Ansible_IoT/Build/templates/bbb_iot_app.service
[Unit]
Description=BBB IoT MQTT Client (bash-wrapped)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/BeagleBoneBlack_IoT
ExecStart=/bin/bash -c 'build/bin/mqtt_mutual'
Restart=no
ProtectHome=no
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib:/usr/lib:/lib"

[Install]
WantedBy=multi-user.target

